# supervisord.conf
[supervisord]
nodaemon=true
[program:node]
command=pm2 start /var/server/src/main.js